# deploy/discord_bot.service
[Unit]
Description=Discord Join Up Bot
After=network.target

[Service]
Type=simple
User=discord_bot
WorkingDirectory=/home/discord_bot/discord_join_up_bot
ExecStart=/home/discord_bot/discord_join_up_bot/run_bot.sh
Restart=on-failure
RestartSec=5

# Optional: redirect stdout/stderr to journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
